# var/lib/dpkg/info/casper.list
/.
/etc
/etc/casper.conf
/etc/init.d
/etc/init.d/casper
/usr
/usr/lib
/usr/lib/casper
/usr/lib/casper/casper-md5check
/usr/share
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/casper
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/casper-functions
/usr/share/initramfs-tools/scripts/casper-bottom
/usr/share/initramfs-tools/scripts/casper-bottom/35fix_language_selector
/usr/share/initramfs-tools/scripts/casper-bottom/34disable_kwallet
/usr/share/initramfs-tools/scripts/casper-bottom/33disable_binary_drivers
/usr/share/initramfs-tools/scripts/casper-bottom/32disable_hibernation
/usr/share/initramfs-tools/scripts/casper-bottom/31disable_update_notifier
/usr/share/initramfs-tools/scripts/casper-bottom/30accessibility
/usr/share/initramfs-tools/scripts/casper-bottom/25configure_init
/usr/share/initramfs-tools/scripts/casper-bottom/24preseed
/usr/share/initramfs-tools/scripts/casper-bottom/23networking
/usr/share/initramfs-tools/scripts/casper-bottom/23etc_modules
/usr/share/initramfs-tools/scripts/casper-bottom/22screensaver
/usr/share/initramfs-tools/scripts/casper-bottom/22gnome_panel_data
/usr/share/initramfs-tools/scripts/casper-bottom/20xconfig
/usr/share/initramfs-tools/scripts/casper-bottom/19keyboard
/usr/share/initramfs-tools/scripts/casper-bottom/18hostname
/usr/share/initramfs-tools/scripts/casper-bottom/15autologin
/usr/share/initramfs-tools/scripts/casper-bottom/14locales
/usr/share/initramfs-tools/scripts/casper-bottom/13swap
/usr/share/initramfs-tools/scripts/casper-bottom/12fstab
/usr/share/initramfs-tools/scripts/casper-bottom/10adduser
/usr/share/initramfs-tools/scripts/casper-bottom/05mountpoints
/usr/share/initramfs-tools/scripts/casper-bottom/01integrity_check
/usr/share/initramfs-tools/scripts/casper-bottom/02etc_casper_conf
/usr/share/initramfs-tools/scripts/casper
/usr/share/initramfs-tools/scripts/casper-helpers
/usr/share/casper
/usr/share/casper/casper-preseed
/usr/share/casper/casper-reconfigure
/usr/share/doc/casper
/usr/share/doc/casper/changelog.gz
/usr/share/doc/casper/copyright
/usr/share/doc/casper/examples
/usr/share/doc/casper/examples/casper.conf
/usr/share/doc/casper/changelog.Debian.gz
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/casper.7.gz
/usr/share/man/man1
/usr/share/man/man1/casper-snapshot.1.gz
/sbin
/sbin/casper-getty
/sbin/casper-login
/sbin/casper-snapshot
/usr/share/man/man7/casper-login.7.gz
/usr/share/man/man7/casper-getty.7.gz
