# var/lib/dpkg/info/cfengine2.list
/.
/etc
/etc/init.d
/etc/init.d/cfengine2
/etc/default
/etc/cfengine
/var
/var/lib
/var/lib/cfengine2
/usr
/usr/sbin
/usr/sbin/cfagent
/usr/sbin/cfkey
/usr/sbin/cfexecd
/usr/sbin/cfservd
/usr/sbin/cfdoc
/usr/sbin/cfenvgraph
/usr/sbin/cfetoolgraph
/usr/sbin/cfrun
/usr/sbin/cfshow
/usr/sbin/cfenvd
/usr/sbin/cfetool
/usr/share
/usr/share/doc/cfengine2
/usr/share/doc/cfengine2/TODO.Debian
/usr/share/doc/cfengine2/copyright
/usr/share/doc/cfengine2/NEWS.Debian.gz
/usr/share/doc/cfengine2/README.Debian
/usr/share/doc/cfengine2/examples
/usr/share/doc/cfengine2/examples/cf.services.example.gz
/usr/share/doc/cfengine2/examples/cf.freebsd.example
/usr/share/doc/cfengine2/examples/cf.ftp.example
/usr/share/doc/cfengine2/examples/cfagent.conf.example
/usr/share/doc/cfengine2/examples/cf.site.example.gz
/usr/share/doc/cfengine2/examples/cf.preconf.example
/usr/share/doc/cfengine2/examples/cf.chflags.example
/usr/share/doc/cfengine2/examples/cf.main.example.gz
/usr/share/doc/cfengine2/examples/cfrun.hosts.example
/usr/share/doc/cfengine2/examples/cf.sun4.example
/usr/share/doc/cfengine2/examples/cfservd.conf.example
/usr/share/doc/cfengine2/examples/cf.linux.example.gz
/usr/share/doc/cfengine2/examples/cfrc.example
/usr/share/doc/cfengine2/examples/update.conf.example
/usr/share/doc/cfengine2/examples/cf.users.example.gz
/usr/share/doc/cfengine2/examples/cf.solaris.example.gz
/usr/share/doc/cfengine2/examples/cfagent.conf-advanced.example
/usr/share/doc/cfengine2/examples/cf.motd.example
/usr/share/doc/cfengine2/examples/cf.groups.example
/usr/share/doc/cfengine2/changelog.Debian.gz
/usr/share/doc/cfengine2/changelog.gz
/usr/share/cfengine2
/usr/share/cfengine2/default
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/cfengine.8.gz
/var/lib/cfengine2/inputs
/var/lib/cfengine2/bin
/usr/share/man/man8/cfagent.8.gz
/usr/share/man/man8/cfenvd.8.gz
/usr/share/man/man8/cfexecd.8.gz
/usr/share/man/man8/cfservd.8.gz
/usr/share/man/man8/cfetool.8.gz
/usr/share/man/man8/cfenvgraph.8.gz
/usr/share/man/man8/cfdoc.8.gz
/usr/share/man/man8/cfkey.8.gz
/usr/share/man/man8/cfrun.8.gz
/usr/share/man/man8/vicf.8.gz
