# var/lib/dpkg/info/dhcp3-server.conffiles
/etc/init.d/dhcp3-server
/etc/dhcp3/dhcpd.conf
