# var/lib/dpkg/info/drqueue.list
/.
/etc
/etc/init.d
/etc/init.d/drqueue-slave
/etc/init.d/drqueue-master
/etc/drqueue
/var
/var/lib
/var/lib/drqueue
/var/lib/drqueue/contrib
/var/lib/drqueue/contrib/sendjob.blender.py
/var/lib/drqueue/etc
/var/lib/drqueue/etc/drqueue
/var/lib/drqueue/etc/drqueue/blender.sg
/var/lib/drqueue/etc/drqueue/pixie.sg
/var/lib/drqueue/etc/drqueue/aqsis.sg
/var/lib/drqueue/etc/drqueue/bmrt.sg
/var/lib/drqueue/etc/drqueue/3delight.sg
/var/lib/drqueue/etc/drqueue/drqman.conf
/var/lib/drqueue/etc/drqueue/maya.sg
/var/lib/drqueue/etc/drqueue/drqman-windows.rc
/var/lib/drqueue/etc/drqueue/mentalray.sg
/var/lib/drqueue/etc/drqueue/slave.conf
/var/lib/drqueue/etc/drqueue/master.conf
/var/lib/drqueue/etc/drqueue/drqman.rc
/var/lib/drqueue/tmp
/var/lib/drqueue/db
/var/lib/drqueue/logs
/var/lib/drqueue/bin
/var/lib/drqueue/bin/blockhost.Linux
/var/lib/drqueue/bin/path2unix.pl
/var/lib/drqueue/bin/drqman
/var/lib/drqueue/bin/master.Linux
/var/lib/drqueue/bin/slave.Linux
/var/lib/drqueue/bin/sendjob.Linux
/var/lib/drqueue/bin/jobfinfo.Linux
/var/lib/drqueue/bin/master
/var/lib/drqueue/bin/cjob.Linux
/var/lib/drqueue/bin/poolcheck.pl
/var/lib/drqueue/bin/sendjob
/var/lib/drqueue/bin/chown_block
/var/lib/drqueue/bin/drqman.Linux
/var/lib/drqueue/bin/requeue.Linux
/var/lib/drqueue/bin/requeue
/var/lib/drqueue/bin/blockhost
/var/lib/drqueue/bin/slasher.pl
/var/lib/drqueue/bin/slave
/var/lib/drqueue/bin/jobfinfo
/usr
/usr/bin
/usr/share
/usr/share/doc/drqueue
/usr/share/doc/drqueue/copyright
/usr/share/doc/drqueue/changelog.Debian.gz
/usr/share/doc/drqueue/changelog.gz
/usr/bin/drqman
/usr/bin/drqueue.master
/usr/bin/drqueue.requeue
/usr/bin/drqueue.jobinfo
/usr/bin/drqueue.blockhost
/usr/bin/drqueue.slave
/usr/bin/drqueue.sendjob
