# var/lib/dpkg/info/drqueue.conffiles
/etc/init.d/drqueue-slave
/etc/init.d/drqueue-master
