# var/lib/dpkg/info/fuse-utils.list
/.
/etc
/etc/init.d
/etc/init.d/fuse
/etc/fuse.conf
/sbin
/sbin/mount.fuse
/usr
/usr/bin
/usr/bin/fusermount
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/fuse-utils
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/fusermount.1.gz
/usr/share/doc/fuse-utils
