# var/lib/dpkg/info/glusterfs-server.list
/.
/var
/var/run
/var/log
/var/log/glusterfs
/etc
/etc/glusterfs
/etc/glusterfs/glusterfs-server.vol.sample
/usr
/usr/sbin
/usr/sbin/glusterfsd
/usr/share
/usr/share/doc/glusterfs-server
/usr/share/doc/glusterfs-server/copyright
/usr/share/doc/glusterfs-server/changelog.Debian.gz
/etc/init.d/glusterfs-server
