# var/lib/dpkg/info/hal.conffiles
/etc/init.d/hal
/etc/default/hal
/etc/dbus-1/system.d/hal.conf
/etc/hal/fdi/policy/preferences.fdi
/etc/udev/rules.d/90-hal.rules
