# var/lib/dpkg/info/htdig.conffiles
/etc/cron.daily/htdig
/etc/htdig/bad_words
/etc/htdig/wrapper.html
/etc/htdig/header.html
/etc/htdig/nomatch.html
/etc/htdig/syntax.html
/etc/htdig/footer.html
/etc/htdig/english.aff
/etc/htdig/long.html
/etc/htdig/htdig.conf
/etc/htdig/synonyms
/etc/htdig/cookies.txt
/etc/htdig/HtFileType-magic.mime
/etc/htdig/english.0
/etc/htdig/mime.types
/etc/htdig/short.html
