# var/lib/dpkg/info/htdig.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/htdig
/etc/htdig
/etc/htdig/bad_words
/etc/htdig/wrapper.html
/etc/htdig/header.html
/etc/htdig/nomatch.html
/etc/htdig/syntax.html
/etc/htdig/footer.html
/etc/htdig/english.aff
/etc/htdig/long.html
/etc/htdig/htdig.conf
/etc/htdig/synonyms
/etc/htdig/cookies.txt
/etc/htdig/HtFileType-magic.mime
/etc/htdig/english.0
/etc/htdig/mime.types
/etc/htdig/short.html
/usr
/usr/include
/usr/include/htdig_db
/usr/include/htdig_db/xa.h
/usr/include/htdig_db/lock_ext.h
/usr/include/htdig_db/xa_ext.h
/usr/include/htdig_db/os_jump.h
/usr/include/htdig_db/txn_auto.h
/usr/include/htdig_db/btree_ext.h
/usr/include/htdig_db/hash_ext.h
/usr/include/htdig_db/mutex_ext.h
/usr/include/htdig_db/crdel_auto.h
/usr/include/htdig_db/btree_auto.h
/usr/include/htdig_db/env_ext.h
/usr/include/htdig_db/common_ext.h
/usr/include/htdig_db/log_auto.h
/usr/include/htdig_db/db_shash.h
/usr/include/htdig_db/db_config.h
/usr/include/htdig_db/db_dispatch.h
/usr/include/htdig_db/qam.h
/usr/include/htdig_db/shqueue.h
/usr/include/htdig_db/lock.h
/usr/include/htdig_db/db_auto.h
/usr/include/htdig_db/region.h
/usr/include/htdig_db/os.h
/usr/include/htdig_db/qam_ext.h
/usr/include/htdig_db/hash.h
/usr/include/htdig_db/hash_auto.h
/usr/include/htdig_db/qam_auto.h
/usr/include/htdig_db/db_page.h
/usr/include/htdig_db/debug.h
/usr/include/htdig_db/mp_ext.h
/usr/include/htdig_db/log.h
/usr/include/htdig_db/queue.h
/usr/include/htdig_db/log_ext.h
/usr/include/htdig_db/mutex.h
/usr/include/htdig_db/txn.h
/usr/include/htdig_db/os_ext.h
/usr/include/htdig_db/db_am.h
/usr/include/htdig_db/btree.h
/usr/include/htdig_db/db.h
/usr/include/htdig_db/db_ext.h
/usr/include/htdig_db/db_join.h
/usr/include/htdig_db/db_int.h
/usr/include/htdig_db/txn_ext.h
/usr/include/htdig_db/mp.h
/usr/include/htdig_db/db_swap.h
/usr/include/htdig
/usr/include/htdig/HtFTP.h
/usr/include/htdig/HtMaxMin.h
/usr/include/htdig/Dictionary.h
/usr/include/htdig/md5.h
/usr/include/htdig/WordDBInfo.h
/usr/include/htdig/Object.h
/usr/include/htdig/Stack.h
/usr/include/htdig/HtRegex.h
/usr/include/htdig/HtRegexReplaceList.h
/usr/include/htdig/HtRandom.h
/usr/include/htdig/good_strtok.h
/usr/include/htdig/List.h
/usr/include/htdig/HtPack.h
/usr/include/htdig/Configuration.h
/usr/include/htdig/htconfig.h
/usr/include/htdig/WordKeyInfo.h
/usr/include/htdig/HtCookieMemJar.h
/usr/include/htdig/HtRegexReplace.h
/usr/include/htdig/Database.h
/usr/include/htdig/DB2_db.h
/usr/include/htdig/WordDBCompress.h
/usr/include/htdig/lib.h
/usr/include/htdig/HtHTTP.h
/usr/include/htdig/HtCodec.h
/usr/include/htdig/ber.h
/usr/include/htdig/StringList.h
/usr/include/htdig/HtCookieInFileJar.h
/usr/include/htdig/WordCursor.h
/usr/include/htdig/myqsort.h
/usr/include/htdig/WordBitCompress.h
/usr/include/htdig/HtCookie.h
/usr/include/htdig/gregex.h
/usr/include/htdig/WordMonitor.h
/usr/include/htdig/Queue.h
/usr/include/htdig/ParsedString.h
/usr/include/htdig/HtVector_String.h
/usr/include/htdig/mhash_md5.h
/usr/include/htdig/WordContext.h
/usr/include/htdig/libdefs.h
/usr/include/htdig/HtHeap.h
/usr/include/htdig/htString.h
/usr/include/htdig/SSLConnection.h
/usr/include/htdig/StringMatch.h
/usr/include/htdig/HtHTTPBasic.h
/usr/include/htdig/WordReference.h
/usr/include/htdig/HtFile.h
/usr/include/htdig/HtCookieJar.h
/usr/include/htdig/HtDateTime.h
/usr/include/htdig/HtWordCodec.h
/usr/include/htdig/WordDBPage.h
/usr/include/htdig/HtVectorGenericCode.h
/usr/include/htdig/HtHTTPSecure.h
/usr/include/htdig/WordType.h
/usr/include/htdig/HtTime.h
/usr/include/htdig/WordStat.h
/usr/include/htdig/HtWordType.h
/usr/include/htdig/HtVector_int.h
/usr/include/htdig/WordRecordInfo.h
/usr/include/htdig/WordKey.h
/usr/include/htdig/HtVectorGeneric.h
/usr/include/htdig/Transport.h
/usr/include/htdig/clib.h
/usr/include/htdig/WordList.h
/usr/include/htdig/WordRecord.h
/usr/include/htdig/HtNNTP.h
/usr/include/htdig/WordCaseIsAStatements.h
/usr/include/htdig/QuotedStringList.h
/usr/include/htdig/HtVector.h
/usr/include/htdig/IntObject.h
/usr/include/htdig/WordDB.h
/usr/include/htdig/Connection.h
/usr/include/htdig/HtRegexList.h
/usr/sbin
/usr/sbin/htdigconfig
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/htdigconfig.8.gz
/usr/share/man/man1
/usr/share/man/man1/htmerge.1.gz
/usr/share/man/man1/htsearch.1.gz
/usr/share/man/man1/htdump.1.gz
/usr/share/man/man1/htnotify.1.gz
/usr/share/man/man1/htload.1.gz
/usr/share/man/man1/htdig.1.gz
/usr/share/man/man1/htfuzzy.1.gz
/usr/share/man/man1/htstat.1.gz
/usr/share/man/man1/htpurge.1.gz
/usr/share/man/man1/rundig.1.gz
/usr/share/man/man1/htdig-pdfparser.1.gz
/usr/share/htdig
/usr/share/htdig/parse_doc.pl
/usr/share/doc/htdig
/usr/share/doc/htdig/README
/usr/share/doc/htdig/changelog.Debian.gz
/usr/share/doc/htdig/README.Debian
/usr/share/doc/htdig/changelog.gz
/usr/share/doc/htdig/examples
/usr/share/doc/htdig/examples/search.html
/usr/share/doc/htdig/copyright
/usr/bin
/usr/bin/htmerge
/usr/bin/htload
/usr/bin/rundig
/usr/bin/htdb_load
/usr/bin/htpurge
/usr/bin/htdb_dump
/usr/bin/htstat
/usr/bin/htnotify
/usr/bin/htfuzzy
/usr/bin/htdump
/usr/bin/htdb_stat
/usr/bin/htdig
/usr/bin/HtFileType
/usr/bin/htdig-pdfparser
/usr/lib
/usr/lib/cgi-bin
/usr/lib/cgi-bin/qtest
/usr/lib/cgi-bin/htsearch
/usr/lib/htdig_db
/usr/lib/htdig_db/libhtdb.la
/usr/lib/htdig_db/libhtdb.a
/usr/lib/htdig_db/libhtdb-3.2.0.so
/usr/lib/htdig
/usr/lib/htdig/libcommon.a
/usr/lib/htdig/libcommon.la
/usr/lib/htdig/libfuzzy-3.2.0.so
/usr/lib/htdig/libhtnet.a
/usr/lib/htdig/libht.a
/usr/lib/htdig/libht-3.2.0.so
/usr/lib/htdig/libfuzzy.a
/usr/lib/htdig/libhtword.a
/usr/lib/htdig/libfuzzy.la
/usr/lib/htdig/libhtnet.la
/usr/lib/htdig/libhtword.la
/usr/lib/htdig/libht.la
/usr/lib/htdig/libcommon-3.2.0.so
/usr/lib/htdig/libhtnet-3.2.0.so
/usr/lib/htdig/libhtword-3.2.0.so
/var
/var/www
/var/www/htdig
/var/www/htdig/htdig.gif
/var/www/htdig/button2.gif
/var/www/htdig/button8.png
/var/www/htdig/star.png
/var/www/htdig/button9.png
/var/www/htdig/button2.png
/var/www/htdig/button4.png
/var/www/htdig/star_blank.png
/var/www/htdig/button10.gif
/var/www/htdig/buttonr.gif
/var/www/htdig/button7.gif
/var/www/htdig/button8.gif
/var/www/htdig/button6.gif
/var/www/htdig/button6.png
/var/www/htdig/button3.gif
/var/www/htdig/button9.gif
/var/www/htdig/button7.png
/var/www/htdig/buttonr.png
/var/www/htdig/buttonl.png
/var/www/htdig/button4.gif
/var/www/htdig/button5.gif
/var/www/htdig/star.gif
/var/www/htdig/star_blank.gif
/var/www/htdig/htdig.png
/var/www/htdig/buttonl.gif
/var/www/htdig/button5.png
/var/www/htdig/button1.png
/var/www/htdig/button1.gif
/var/www/htdig/button3.png
/var/www/htdig/button10.png
/var/lib
/var/lib/htdig
/usr/lib/htdig_db/libhtdb.so
/usr/lib/htdig/libhtnet.so
/usr/lib/htdig/libfuzzy.so
/usr/lib/htdig/libht.so
/usr/lib/htdig/libhtword.so
/usr/lib/htdig/libcommon.so
