# var/lib/dpkg/info/kbd.conffiles
/etc/init.d/console-screen.kbd.sh
/etc/kbd/config
/etc/kbd/remap
