# var/lib/dpkg/info/libdevmapper1.02.list
/.
/etc
/etc/init.d
/etc/init.d/libdevmapper1.02
/usr
/usr/share
/usr/share/doc/libdevmapper1.02
/usr/share/doc/libdevmapper1.02/changelog.Debian.gz
/usr/share/doc/libdevmapper1.02/copyright
/usr/share/doc/libdevmapper1.02/changelog.gz
/lib
/lib/libdevmapper.so.1.02
