# var/lib/dpkg/info/scrollkeeper.conffiles
/etc/cron.monthly/scrollkeeper
/etc/logrotate.d/scrollkeeper
/etc/scrollkeeper.conf
