# var/lib/dpkg/info/splashy.list
/.
/sbin
/sbin/splashy_update
/sbin/splashy_config
/sbin/splashy
/etc
/etc/console-tools
/etc/console-tools/config.d
/etc/console-tools/config.d/splashy
/etc/init.d
/etc/init.d/splashy
/etc/splashy
/etc/splashy/config.xml
/etc/default
/etc/default/splashy
/etc/kbd
/etc/kbd/config.d
/etc/lsb-base-logging.sh
/usr
/usr/share
/usr/share/initramfs-tools
/usr/share/initramfs-tools/modules.d
/usr/share/initramfs-tools/modules.d/splashy
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-top
/usr/share/initramfs-tools/scripts/init-top/splashy
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/libsplashy
/usr/share/initramfs-tools/hooks/splashy
/usr/share/locale
/usr/share/locale/sk/LC_MESSAGES/splashy.mo
/usr/share/locale/de/LC_MESSAGES/splashy.mo
/usr/share/locale/pt/LC_MESSAGES/splashy.mo
/usr/share/locale/fr/LC_MESSAGES/splashy.mo
/usr/share/locale/es/LC_MESSAGES/splashy.mo
/usr/share/locale/hu/LC_MESSAGES/splashy.mo
/usr/share/splashy
/usr/share/splashy/schemas
/usr/share/splashy/schemas/config.xsd
/usr/share/splashy/schemas/theme.xsl
/usr/share/splashy/schemas/theme.xsd
/usr/share/splashy/themes
/usr/share/splashy/themes/default
/usr/share/splashy/themes/default/error.png
/usr/share/splashy/themes/default/suspend.png
/usr/share/splashy/themes/default/FreeSans.ttf
/usr/share/splashy/themes/default/background.png
/usr/share/splashy/themes/default/theme.xml
/usr/share/doc/splashy
/usr/share/doc/splashy/TODO
/usr/share/doc/splashy/changelog.gz
/usr/share/doc/splashy/copyright
/usr/share/doc/splashy/changelog.Debian.gz
/usr/share/doc/splashy/README.gz
/usr/share/doc/splashy/README.Debian
/usr/share/doc/splashy/AUTHORS
/usr/share/doc/splashy/NEWS.gz
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/splashy-theme.xml.5.gz
/usr/share/man/man5/splashy-config.xml.5.gz
/usr/share/man/man1
/usr/share/man/man1/splashy_update.1.gz
/usr/share/man/man1/splashy_config.1.gz
/usr/share/man/man1/splashy.1.gz
/sbin/splashy_chvt
/etc/splashy/themes
/etc/kbd/config.d/splashy
/usr/share/man/man1/splashy_chvt.1.gz
