# var/lib/dpkg/info/splashy.conffiles
/etc/console-tools/config.d/splashy
/etc/init.d/splashy
/etc/splashy/config.xml
/etc/default/splashy
/etc/lsb-base-logging.sh
