# var/lib/dpkg/info/systemimager-server.conffiles
/etc/systemimager/systemimager.conf
/etc/systemimager/autoinstallscript.template
/etc/systemimager/imagemanip.conf
/etc/systemimager/imagemanip.perm
/etc/systemimager/rsync_stubs/10header
/etc/systemimager/rsync_stubs/99local
/etc/systemimager/rsync_stubs/README
/etc/systemimager/pxelinux.cfg/message.txt
/etc/systemimager/pxelinux.cfg/syslinux.cfg
/etc/systemimager/pxelinux.cfg/syslinux.cfg.localboot
/etc/systemimager/pxelinux.cfg/default
/etc/init.d/systemimager-server
/etc/init.d/netbootmond
