# var/lib/dpkg/info/systemimager-boot-i386-standard.list
/.
/usr
/usr/share
/usr/share/systemimager
/usr/share/systemimager/boot
/usr/share/systemimager/boot/i386
/usr/share/systemimager/boot/i386/standard
/usr/share/systemimager/boot/i386/standard/kernel
/usr/share/systemimager/boot/i386/standard/config
/usr/share/systemimager/boot/i386/standard/initrd.img
/usr/share/systemimager/boot/i386/standard/boel_binaries.tar.gz
/usr/share/systemimager/boot/i386/standard/initrd_template
/usr/share/systemimager/boot/i386/standard/initrd_template/bin
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/busybox
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/ash
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/chmod
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/cat
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/udp-receiver
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/chown
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/cp
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/dd
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/df
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/echo
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/grep
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/kill
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/ln
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/ls
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/mkdir
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/mknod
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/more
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/mount
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/mv
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/ping
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/ps
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/rm
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/rsync
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/sed
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/sh
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/shutdown
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/sleep
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/tar
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/umount
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/uname
/usr/share/systemimager/boot/i386/standard/initrd_template/bin/vi
/usr/share/systemimager/boot/i386/standard/initrd_template/dev
/usr/share/systemimager/boot/i386/standard/initrd_template/dev/console
/usr/share/systemimager/boot/i386/standard/initrd_template/dev/null
/usr/share/systemimager/boot/i386/standard/initrd_template/etc
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/init.d
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/init.d/functions
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/init.d/rcS
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/dhclient.conf
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/devfsd.conf
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/dhclient-script.systemimager-prefix
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/dhclient-script.debian-dist
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/nsswitch.conf
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/group
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/hosts
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/inittab
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/issue
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/mtab
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/passwd
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/profile
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/services
/usr/share/systemimager/boot/i386/standard/initrd_template/etc/termcap
/usr/share/systemimager/boot/i386/standard/initrd_template/lib
/usr/share/systemimager/boot/i386/standard/initrd_template/lib/libnss_dns-2.3.6.so
/usr/share/systemimager/boot/i386/standard/initrd_template/lib/ld-linux.so.1
/usr/share/systemimager/boot/i386/standard/initrd_template/lib/ld-linux.so.2
/usr/share/systemimager/boot/i386/standard/initrd_template/lib/ld.so
/usr/share/systemimager/boot/i386/standard/initrd_template/lib/ld.so.1.9.11
/usr/share/systemimager/boot/i386/standard/initrd_template/lib/libc.so.6
/usr/share/systemimager/boot/i386/standard/initrd_template/lib/libdl.so.2
/usr/share/systemimager/boot/i386/standard/initrd_template/lib/libpthread.so.0
/usr/share/systemimager/boot/i386/standard/initrd_template/lib/libresolv.so.2
/usr/share/systemimager/boot/i386/standard/initrd_template/my_modules
/usr/share/systemimager/boot/i386/standard/initrd_template/my_modules/INSMOD_COMMANDS
/usr/share/systemimager/boot/i386/standard/initrd_template/new_root
/usr/share/systemimager/boot/i386/standard/initrd_template/proc
/usr/share/systemimager/boot/i386/standard/initrd_template/root
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/pivot_root
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/devfsd
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/dhclient
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/halt
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/ifconfig
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/init
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/insmod
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/klogd
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/lsmod
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/modprobe
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/poweroff
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/reboot
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/route
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/swapoff
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/swapon
/usr/share/systemimager/boot/i386/standard/initrd_template/sbin/syslogd
/usr/share/systemimager/boot/i386/standard/initrd_template/tmp
/usr/share/systemimager/boot/i386/standard/initrd_template/usr
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/killall
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/[
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/cut
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/expr
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/logger
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/nc
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/nslookup
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/sort
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/test
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/tr
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/uniq
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/bin/wget
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/sbin
/usr/share/systemimager/boot/i386/standard/initrd_template/usr/sbin/chroot
/usr/share/systemimager/boot/i386/standard/initrd_template/var
/usr/share/systemimager/boot/i386/standard/initrd_template/var/run
/usr/share/systemimager/boot/i386/standard/initrd_template/linuxrc
/usr/share/doc/systemimager-boot-i386-standard
/usr/share/doc/systemimager-boot-i386-standard/copyright
/usr/share/doc/systemimager-boot-i386-standard/changelog.Debian.gz
/usr/share/systemimager/boot/i386/standard/initrd_template/lib/libnss_dns.so.2
