# var/lib/dpkg/info/systemimager-server.list
/.
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/si_addclients.8.gz
/usr/share/man/man8/si_getimage.8.gz
/usr/share/man/man8/si_mkautoinstallcd.8.gz
/usr/share/man/man8/si_mkautoinstalldiskette.8.gz
/usr/share/man/man8/si_mkautoinstallscript.8.gz
/usr/share/man/man8/si_mkdhcpserver.8.gz
/usr/share/man/man8/si_mkdhcpstatic.8.gz
/usr/share/man/man8/si_cpimage.8.gz
/usr/share/man/man8/si_mvimage.8.gz
/usr/share/man/man8/si_rmimage.8.gz
/usr/share/man/man8/si_mkbootserver.8.gz
/usr/share/man/man8/si_mkrsyncd_conf.8.gz
/usr/share/man/man8/si_pushupdate.8.gz
/usr/share/man/man8/si_mkclientnetboot.8.gz
/usr/share/man/man5
/usr/share/man/man5/systemimager.conf.5.gz
/usr/share/man/man5/systemimager-rsyncd.conf.5.gz
/usr/share/systemimager
/usr/share/systemimager/boot
/usr/share/systemimager/boot/ia64
/usr/share/systemimager/boot/ia64/standard
/usr/share/doc/systemimager-server
/usr/share/doc/systemimager-server/README.Debian
/usr/share/doc/systemimager-server/copyright
/usr/share/doc/systemimager-server/examples
/usr/share/doc/systemimager-server/examples/elilo-ia64.conf
/usr/share/doc/systemimager-server/examples/post-install
/usr/share/doc/systemimager-server/examples/post-install/99all.xcat_updateflag
/usr/share/doc/systemimager-server/examples/post-install/99ia64-compute.set_efi_boot_order
/usr/share/doc/systemimager-server/examples/pxelinux.cfg
/usr/share/doc/systemimager-server/examples/syslinux.cfg
/usr/share/doc/systemimager-server/examples/user_contrib
/usr/share/doc/systemimager-server/examples/user_contrib/si-client
/usr/share/doc/systemimager-server/examples/user_contrib/README
/usr/share/doc/systemimager-server/examples/yaboot.conf
/usr/share/doc/systemimager-server/examples/local.cfg.gz
/usr/share/doc/systemimager-server/examples/screencaptured_multicast_install.txt.gz
/usr/share/doc/systemimager-server/changelog.Debian.gz
/usr/lib
/usr/lib/systemimager
/usr/lib/systemimager/perl
/usr/lib/systemimager/perl/SystemImager
/usr/lib/systemimager/perl/SystemImager/Server.pm
/usr/bin
/usr/bin/si_mkautoinstallcd
/usr/bin/si_mkautoinstalldiskette
/usr/sbin
/usr/sbin/si_addclients
/usr/sbin/si_cpimage
/usr/sbin/si_getimage
/usr/sbin/si_mkdhcpserver
/usr/sbin/si_mkdhcpstatic
/usr/sbin/si_mkautoinstallscript
/usr/sbin/si_mkbootserver
/usr/sbin/si_mvimage
/usr/sbin/si_pushupdate
/usr/sbin/si_rmimage
/usr/sbin/si_mkrsyncd_conf
/usr/sbin/si_mkclientnetboot
/usr/sbin/si_netbootmond
/usr/sbin/si_imagemanip
/usr/sbin/si_mkbootpackage
/usr/sbin/si_monitor
/usr/sbin/si_monitortk
/etc
/etc/systemimager
/etc/systemimager/systemimager.conf
/etc/systemimager/autoinstallscript.template
/etc/systemimager/imagemanip.conf
/etc/systemimager/imagemanip.perm
/etc/systemimager/rsync_stubs
/etc/systemimager/rsync_stubs/10header
/etc/systemimager/rsync_stubs/99local
/etc/systemimager/rsync_stubs/README
/etc/systemimager/pxelinux.cfg
/etc/systemimager/pxelinux.cfg/message.txt
/etc/systemimager/pxelinux.cfg/syslinux.cfg
/etc/systemimager/pxelinux.cfg/syslinux.cfg.localboot
/etc/systemimager/pxelinux.cfg/default
/etc/init.d
/etc/init.d/systemimager-server
/etc/init.d/netbootmond
/var
/var/log
/var/log/systemimager
/var/lock
/var/lock/systemimager
/var/lib
/var/lib/systemimager
/var/lib/systemimager/scripts
/var/lib/systemimager/scripts/pre-install
/var/lib/systemimager/scripts/pre-install/99all.harmless_example_script
/var/lib/systemimager/scripts/post-install
/var/lib/systemimager/scripts/post-install/99all.harmless_example_script
/var/lib/systemimager/scripts/post-install/README
/var/lib/systemimager/overrides
/var/lib/systemimager/overrides/README
/var/lib/systemimager/images
/var/lib/systemimager/images/README
/var/lib/systemimager/images/CUIDADO
/var/lib/systemimager/images/ACHTUNG
/var/lib/systemimager/images/DO_NOT_TOUCH_THESE_DIRECTORIES
/var/state
/var/state/systemimager
/var/state/systemimager/flamethrower
