# var/lib/dpkg/info/udev.list
/.
/usr
/usr/share
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-premount
/usr/share/initramfs-tools/scripts/init-premount/udev
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/udev
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/udev
/usr/share/bug
/usr/share/bug/udev
/usr/share/bug/udev/presubj
/usr/share/bug/udev/script
/usr/share/doc/udev
/usr/share/doc/udev/TODO
/usr/share/doc/udev/README.vol_id
/usr/share/doc/udev/changelog.gz
/usr/share/doc/udev/copyright
/usr/share/doc/udev/README
/usr/share/doc/udev/examples
/usr/share/doc/udev/examples/udev.vim
/usr/share/doc/udev/NEWS.Debian.gz
/usr/share/doc/udev/RELEASE-NOTES.gz
/usr/share/doc/udev/README.Debian.gz
/usr/share/doc/udev/changelog.Debian.gz
/usr/share/doc/udev/writing_udev_rules
/usr/share/doc/udev/writing_udev_rules/index.html
/usr/share/doc/udev/FAQ.gz
/usr/share/doc-base
/usr/share/doc-base/udev
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/udev.7.gz
/usr/share/man/man8
/usr/share/man/man8/scsi_id.8.gz
/usr/share/man/man8/udevadm.8.gz
/usr/share/man/man8/cdrom_id.8.gz
/usr/share/man/man8/edd_id.8.gz
/usr/share/man/man8/ata_id.8.gz
/usr/share/man/man8/vol_id.8.gz
/usr/share/man/man8/udevd.8.gz
/usr/share/man/man8/path_id.8.gz
/usr/lib
/usr/lib/hotplug
/usr/lib/hotplug/firmware
/usr/bin
/lib
/lib/firmware
/lib/udev
/lib/udev/usb_id
/lib/udev/path_id
/lib/udev/devices
/lib/udev/firmware.agent
/lib/udev/net.agent
/lib/udev/collect
/lib/udev/vol_id
/lib/udev/hotplug.functions
/lib/udev/scsi-re-add
/lib/udev/write_cd_rules
/lib/udev/ata_id
/lib/udev/cdrom_id
/lib/udev/dsl-modem.agent
/lib/udev/edd_id
/lib/udev/write_net_rules
/lib/udev/logger.agent
/lib/udev/scsi_id
/lib/udev/fstab_import
/lib/udev/create_static_nodes
/etc
/etc/init.d
/etc/init.d/udev-mtab
/etc/init.d/udev
/etc/modprobe.d
/etc/modprobe.d/blacklist
/etc/modprobe.d/display_class
/etc/modprobe.d/pnp-hotplug
/etc/scsi_id.config
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/50-udev.rules
/etc/udev/rules.d/91-permissions.rules
/etc/udev/rules.d/60-persistent-v4l.rules
/etc/udev/rules.d/75-persistent-net-generator.rules
/etc/udev/rules.d/75-cd-aliases-generator.rules
/etc/udev/rules.d/80-drivers.rules
/etc/udev/rules.d/95-late.rules
/etc/udev/rules.d/60-persistent-input.rules
/etc/udev/rules.d/60-persistent-storage.rules
/etc/udev/rules.d/60-persistent-storage-tape.rules
/etc/udev/links.conf
/etc/udev/udev.conf
/etc/udev/.dev
/sbin
/sbin/udevd
/sbin/udevadm
/usr/share/man/man8/udevinfo.8.gz
/usr/share/man/man8/udevsettle.8.gz
/usr/bin/udevinfo
/sbin/vol_id
/sbin/udevsettle
