# var/lib/dpkg/info/xen-utils-common.list
/.
/etc
/etc/init.d
/etc/init.d/xend
/etc/init.d/xendomains
/etc/default
/etc/default/xendomains
/etc/udev
/etc/udev/xen-backend.rules
/etc/xen
/etc/xen/xend-config-xenapi.sxp
/etc/xen/xend-config.sxp
/etc/xen/xend-pci-quirks.sxp
/etc/xen/xend-pci-permissive.sxp
/etc/xen/scripts
/etc/xen/scripts/vtpm-delete
/etc/xen/scripts/qemu-ifup
/etc/xen/scripts/vtpm-common.sh
/etc/xen/scripts/logging.sh
/etc/xen/scripts/network-route
/etc/xen/scripts/external-device-migrate
/etc/xen/scripts/xen-network-common.sh
/etc/xen/scripts/block-enbd
/etc/xen/scripts/vif-route
/etc/xen/scripts/vif-nat
/etc/xen/scripts/vif-common.sh
/etc/xen/scripts/vtpm
/etc/xen/scripts/network-bridge
/etc/xen/scripts/vtpm-impl
/etc/xen/scripts/block-nbd
/etc/xen/scripts/locking.sh
/etc/xen/scripts/block-common.sh
/etc/xen/scripts/network-dummy
/etc/xen/scripts/blktap
/etc/xen/scripts/block
/etc/xen/scripts/vtpm-hotplug-common.sh
/etc/xen/scripts/network-nat
/etc/xen/scripts/vif-bridge
/etc/xen/scripts/vtpm-migration.sh
/etc/xen/scripts/xen-script-common.sh
/etc/xen/scripts/xen-hotplug-cleanup
/etc/xen/scripts/xen-hotplug-common.sh
/usr
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/xm.1.gz
/usr/share/man/man5
/usr/share/man/man5/xend-config.sxp.5.gz
/usr/share/man/man5/xmdomain.cfg.5.gz
/usr/share/doc/xen-utils-common
/usr/share/doc/xen-utils-common/changelog.Debian.gz
/usr/share/doc/xen-utils-common/copyright
/usr/share/doc/xen-utils-common/examples
/usr/share/doc/xen-utils-common/examples/xm-config.xml
/usr/share/doc/xen-utils-common/examples/xmexample2.gz
/usr/share/doc/xen-utils-common/examples/xmexample.hvm.gz
/usr/share/doc/xen-utils-common/examples/xmexample.vti.gz
/usr/share/doc/xen-utils-common/examples/xm-config-xenapi.xml
/usr/share/doc/xen-utils-common/examples/xmexample1.gz
/usr/sbin
/usr/lib
/usr/lib/xen-common
/usr/lib/xen-common/bin
/usr/lib/xen-common/bin/xen-utils-wrapper
/usr/lib/xen-common/bin/xen-utils-version
/usr/sbin/xentrace_setsize
/usr/sbin/cpuperf-perfcntr
/usr/sbin/xend
/usr/sbin/xentrace
/usr/sbin/cpuperf-xen
/usr/sbin/xenperf
/usr/sbin/xm
/usr/sbin/xentrace_format
/usr/sbin/xentop
/usr/sbin/xentrace_setmask
